# etc/supervisor/conf.d/orocrm-websocket.conf
[program:orocrm_web_socket]
command=${OROCRM_ROOT_DIR}/app/console clank:server --env=${SYMFONY_ENV}
numprocs=1
autostart=true
startsecs=0
user=${APP_USER}
redirect_stderr=true