# etc/cron.d/orocrm.crontab
*/1 * * * * ${APP_USER} ${OROCRM_ROOT_DIR}/app/console oro:cron --env=${SYMFONY_ENV}